# 02_recursos/contactanos.socket
[Unit]
Description=Socket de la API Contactanos

[Socket]
ListenStream=/run/contactanos.sock

[Install]
WantedBy=sockets.target
